# launchd/com.user.hhmonitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.hhmonitor</string>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/hh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/hh/service-run</string>
    <string>--settings</string>
    <string>config/settings.json</string>
    <string>--mode</string>
    <string>fast</string>
    <string>--max-pages</string>
    <string>10</string>
    <string>--max-age-days</string>
    <string>30</string>
    <string>--export</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/hh/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/hh/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
